# templates/contentserver/jms.service.erb
[Unit]
Description=MethodServer
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/nohup <%=@documentum%>/wildfly9.0.1/server/startMethodServer.sh &>/dev/null &
ExecStop=<%=@documentum%>/shared/wildfly9.0.1/server/stopMethodServer.sh
RemainAfterExit=yes
EnvironmentFile=/etc/default/dctm
User=<%=@dctm_owner%>

[Install]
WantedBy=multi-user.target
